# templates/motd/update_motd.service.j2
[Unit]
Description=update motd
After=sshdgenkeys.service

[Service]
ExecStart=
ExecStart=sh /usr/bin/update_motd.sh
Restart=always

[Install]
WantedBy=multi-user.target
